# apps/backend-main/docker/supervisord/supervisord_programs/nginx.conf
[program:nginx]
directory=/app
command=/usr/sbin/nginx
stdout_logfile=./.logs/nginx.output.log
stdout_logfile_maxbytes=50MB
stderr_logfile=./.logs/nginx.error.log
stderr_logfile_maxbytes=50MB
